# .config/systemd/user/ncc-record.timer
[Unit]
Description=Trigger the service record.service
After=autofs.service

[Timer]
OnCalendar=Sun 09:29:30 UTC

[Install]
WantedBy=timers.target
