# linux/ubuntu/common/usr/local/lib/systemd/system/ncc-record.timer
[Unit]
Description=Trigger the service record.service

[Timer]
OnCalendar=Sun *-*-* 08:29:30

[Install]
WantedBy=timers.target
