# .config/autostart/conky.desktop
[Desktop Entry]
Comment=
Comment[en_IN]=
Exec=sh "sleep 30;/home/sebastien/myScripts/shl/conky_restart.sh"
Name=Conky
Name[en_IN]=Conky
NoDisplay=false
Type=Application
X-GNOME-Autostart-enabled=true
