# .config/autostart/fluxgui.desktop
[Desktop Entry]
Exec=fluxgui
Icon=fluxgui
Name=f.lux indicator applet
Type=Application
X-GNOME-Autostart-enabled=true
