# .config/autostart/Play_With_MPV-server.desktop
[Desktop Entry]
Categories=AudioVideo;Audio;Video;Player;TV
Exec=/usr/local/bin/play-with-mpv
Icon=mpv
Name=Play With MPV (server)
Type=Application
