# linux/systemd_examples/etc/systemd/system/haproxy.service.d/override.conf
[Service]
Environment=-
Environment=CONFIG=/etc/haproxy

