# linux/systemd_examples/etc/systemd/system/haproxy.service.d/override.conf
[Service]
EnvironmentFile=
Environment=
Environment=CONFIG=/etc/haproxy
